# dorfleaks.timer
[Unit]
Description=create dorfleaks.json every day

[Timer]
OnCalendar=daily
Persistent=True

[Install]
WantedBy=timers.target
